# supervisor-kiran.conf
[program:kiran-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/kiranindustries/artisan queue:work --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=Administrator
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/kiranindustries/storage/logs/queue-worker.log

[program:kiran-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/kiranindustries/artisan schedule:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=Administrator
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/kiranindustries/storage/logs/scheduler.log
; Scheduled tasks (via schedule:work):
;   activity-logs:cleanup — daily 03:00, deletes logs older than 2 months
;   reports:cleanup       — daily
;   biometric:sync-all    — every 30 minutes
